# scripts/supervisord.ini
[supervisord]
nodaemon=true
user=root

[program:uwsgi]
command=/usr/sbin/uwsgi /etc/uwsgi/uwsgi.ini
stdout_logfile=/var/log/messages
stdout_logfile_maxbytes=100000
stderr_logfile=/var/log/messages
stderr_logfile_maxbytes=100000
stdout_logfile_backups=0
stderr_logfile_backups=0
